# restic-backup.timer
[Unit]
Description=Restic backup timer (daily at 2:00 AM)
Requires=restic-backup.service

[Timer]
# Run daily at 2:00 AM
OnCalendar=*-*-* 02:00:00

# If system was off during scheduled time, run on next boot
Persistent=true

# Add randomization to avoid conflicts (up to 15 minutes)
RandomizedDelaySec=15min

[Install]
WantedBy=timers.target
